# recipes-framework/packet-forwarder/packet-forwarder/lorawan-gateway-m4.service
[Unit]
Description=LoRaWan Gateway

[Service]
WorkingDirectory=/usr/local/lorawan-gateway
ExecStart=/usr/local/lorawan-gateway/launch_m4.sh
SyslogIdentifier=stm32mp1-lorawan-m4
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
